# kanban/runbooks/fleet-conductor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Copy to ~/Library/LaunchAgents/com.migx.fleet-conductor.plist
  Edit WorkingDirectory and paths, then:
    launchctl load ~/Library/LaunchAgents/com.migx.fleet-conductor.plist
  Unload:
    launchctl unload ~/Library/LaunchAgents/com.migx.fleet-conductor.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.migx.fleet-conductor</string>
  <key>WorkingDirectory</key>
  <string>/Users/gudjon/code/migx</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/gudjon/code/migx/kanban/scripts/migx-fleet-conductor.py</string>
    <string>--nudge-file</string>
    <string>--drain-codex</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/migx-fleet-conductor.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/migx-fleet-conductor.err</string>
</dict>
</plist>
